# jug-server/jira-data-loader.cron
*/5 * * * * loader /opt/loader/bin/recent
